# templates/project.service.erb
#############################################################################
# Example: VsMyprojects/examples/TestApplication/TestApplication.service
#############################################################################
[Unit]
Description=<%= @description %>

[Service]
WorkingDirectory=/srv/<%= @projectName %>
ExecStart=/srv/<%= @projectName %>/<%= @application %>
SyslogIdentifier=<%= @projectName %>
User=<%= @sdkUser %>
Restart=always
RestartSec=5

KillSignal=SIGINT
Environment=ASPNETCORE_ENVIRONMENT=Development
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
#Environment=DOTNET_ROOT=/home/<%= @sdkUser %>/.dotnet/dotnet
Environment=ASPNETCORE_URLS=<%= @aspnetCoreUrls %>

<% if @ssl -%>
Environment=ASPNETCORE_Kestrel__Certificates__Default__Path=/etc/pki/tls/certs/<%= @sslHost %>.crt
Environment=ASPNETCORE_Kestrel__Certificates__Default__KeyPath=/etc/pki/tls/private/<%= @sslHost %>.key
<% end -%>

[Install]
WantedBy=multi-user.target
